# etc/profile.d/00-var-order.sh
# Export VAR_ORDER before trying to manipulate it later.

export VAR_ORDER
